# debian/template/etc/init/ryft-server-d.conf
description "ryft-server-d"
author "RYFT"

### run as ryftuser instead of root
# setuid ryftuser
# setgid ryftuser

start on runlevel [2345] and started rhfsd
stop on shutdown

pre-start script
	echo `date "+%F %T"` "pre-start"

	# wait for partition
	waitSeconds=15
	partition="/ryftone"

	begin=`date +%s`
	end=$((begin+waitSeconds))
	partitionMounted=$(df -h | grep $partition | wc -l)
	while [ "$partitionMounted" = 0 -a `date +%s` -lt $end ] 
	do
		# retry to find partition
		sleep 3 
		partitionMounted=$(df -h | grep $partition | wc -l)
	done
	if [ "$partitionMounted" = 0 ] 
	then
		# unable to find partition
		logger "ryft-server-d pre-start unable to find $partition partition"
		echo `date "+%F %T"` "pre-start unable to find $partition partition"
		exit 1;
	fi
end script

post-stop script
	echo `date "+%F %T"` "post-stop"
end script

script
	DAEMON_PATH=/usr/bin
	
	DAEMON=ryft-server
	CONFIGFILE="/etc/ryft-server.conf"
	
	NAME=ryft-server-d
	DESC="Ryft REST API"
	PIDFILE=/var/run/$NAME.pid
	SCRIPTNAME=/etc/init.d/$NAME
	
	export TERM=xterm-256color
	export LC_ALL=en_US.UTF-8
	export PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
	export LANG=en_US.UTF-8
	export HOME=/home/ryftuser
	export LC_CTYPE=UTF-8
	export SHELL=/bin/bash
	export MAIL=/var/mail/root
	export LOGNAME=root
	export USER=root
	export USERNAME=root
	export SUDO_COMMAND=/usr/bin/printenv
	export SUDO_USER=ryftuser
	export SUDO_UID=1001
	export SUDO_GID=1001
	
	printf "\n%-50s" "Starting $NAME..."
	cd $DAEMON_PATH
	
	if [ -f $CONFIGFILE ]; then

		#exec $DAEMON --config=$CONFIGFILE
		exec $DAEMON --config=$CONFIGFILE >/var/log/ryft/server.log 2>&1
	else
		echo
		echo There is no config $CONFIGFILE was found in /etc
		echo Default settings will be used

		#exec $DAEMON
		exec $DAEMON >/var/log/ryft/server.log 2>&1
	fi
end script
