# debian/template/DEBIAN/postinst
#!/bin/bash

NAME=ryft-server-d

# force destination uid gid if needed
target_uid=ryftuser
target_gid=ryftuser
function fixUIDGID {
	if [[ `ls -ald $1|awk '{print \$3}'` != $target_uid ]]
	then
		chown $target_uid $1
	fi
	if [[ `ls -ald $1|awk '{print \$4}'` != $target_gid ]]
	then
		chown :$target_gid $1
	fi
}

# executables
fixUIDGID /usr/bin/ryft_ps
fixUIDGID /usr/bin/ryft_rdf
fixUIDGID /usr/bin/ryft_status
fixUIDGID /usr/bin/ryft-server
fixUIDGID /usr/bin/ryft-server-aggs
# remove ryftrest and ryft-versions in future
fixUIDGID /usr/bin/ryftrest
fixUIDGID /usr/bin/blnxrest
fixUIDGID /usr/bin/ryftutil
fixUIDGID /usr/bin/ryft-versions
fixUIDGID /usr/bin/blnx-versions
fixUIDGID /usr/bin/ryftx_wrapper

fixUIDGID /usr/bin/cpc

# configurations
chown "$target_uid:$target_gid" -f /etc/ryft-server.conf
chown "$target_uid:$target_gid" -f /etc/ryft-users.yaml

# other
chown "$target_uid:$target_gid" -f /var/log/ryft
chown "$target_uid:$target_gid" -f -R /var/ryft
chown "$target_uid:$target_gid" -f -R /home/ryftuser/.blacklynx/bl-server.orig
if [ ! -f /home/ryftuser/.blacklynx/bl-server.cfg ]; then
	cp /home/ryftuser/.blacklynx/bl-server.orig /home/ryftuser/.blacklynx/bl-server.cfg
fi

# backward compatibility
chmod 0755 -f /ryftone/.rest-*
chown "$target_uid:$target_gid" -f -R /ryftone/.rest-*
chown "$target_uid:$target_gid" -f /var/log/blacklynx/server.log

# pull the Docker images (check the docker.images configuration options)
[ -z "$RYFT_DO_NOT_PULL_DOCKER_IMAGES" ] && \
(( docker pull alpine:latest \
&& docker pull ubuntu:16.04 \
&& docker pull python:2.7 \
&& docker build -t ryft/alpine-tshark /usr/share/ryft/docker/tshark ) \
|| echo "failed to pull some Docker images, the /run endpoint may fail" )

echo $NAME starting...
if [ -z `which systemctl` ]; then
	# yes upstart
	chown "$target_uid:$target_gid" -f /etc/init/ryft-server-d.conf
	chown "$target_uid:$target_gid" -f /etc/init/rhfsdTrigger.conf

	# no systemd
	rm -f /etc/systemd/system/ryft-server-d.service

	# start service (if no $RYFT_DO_NOT_START defined)
	[ -z "$RYFT_DO_NOT_START" ] && service $NAME start || true
else
	# yes systemd
	chown root:root -f /etc/systemd/system/ryft-server-d.service
	systemctl enable ryft-server-d

	# no upstart
	rm -f /etc/init/ryft-server-d.conf
	rm -f /etc/init/rhfsdTrigger.conf

	# start service (if no $RYFT_DO_NOT_START defined)
	[ -z "$RYFT_DO_NOT_START" ] && systemctl start $NAME || true
fi
