# debian/template/DEBIAN/postinst
#!/bin/bash

echo ryft-server-d installing for autostart...
update-rc.d ryft-server-d defaults

echo ryft-server-d starting...
service ryft-server-d start
