# debian/template/DEBIAN/prerm
#!/bin/bash

NAME=ryft-server-d

echo $NAME uninstalling from autostart...
update-rc.d -f $NAME remove

echo $NAME stopping...
service $NAME stop

# have to return success
exit 0
